# docker/crontab
# Crontab embarquée pour le mode `docker run plex-scripts cron`
# Remplace les timers systemd. Les logs sont envoyés sur stdout (→ docker logs).
#
# Format : m h dom mon dow  commande
# TZ est fixé par la variable d'env du conteneur (défaut Europe/Paris).

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
PYTHONUNBUFFERED=1

# MODE JOURNÉE — 1 exécution par heure (08h à 22h), commandes décalées
# 00' : Sync ratings Plex → ID3
0 8-22 * * *   cd /app && python3 ratings/plex_rating_sync_complete.py >> /data/logs/ratings.log 2>&1

# 15' : Workflow 1★/2★
15 8-22 * * *  cd /app && ./workflows/plex_daily_workflow.sh >> /data/logs/daily.log 2>&1

# 30' : Génération playlists Plexamp
30 8-22 * * *  cd /app && ./playlists/generate_plexamp_playlists.sh --refresh >> /data/logs/playlists.log 2>&1

# 45' : Export des playlists en M3U
45 8-22 * * *  cd /app && python3 playlists/export_plex_playlist.py --all --dest /playlists >> /data/logs/export.log 2>&1
