# systemd/plex-playlists-mybook-sync.service
[Unit]
Description=Sync Plex playlists to MyBook export directory
After=snap.plexmediaserver.plexmediaserver.service
RequiresMountsFor=/mnt/MyBook

[Service]
Type=oneshot
User=paulceline
ExecStartPre=/bin/bash -c 'mountpoint -q /mnt/MyBook || { echo "MyBook non monté"; exit 1; }'
ExecStart=/home/paulceline/scripts/playlists/sync_plex_playlists_to_mybook.sh

[Install]
WantedBy=multi-user.target
